# deploy/runyard-runner.service
[Unit]
Description=Runyard Runner
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/runyard
EnvironmentFile=/etc/runyard/runner.env
Environment=RUNYARD_HUB_URL=https://hub.example.com
Environment="SMITHERS_RUNNER_NAME=vps-runner-1"
Environment=SMITHERS_RUNNER_LOCATION=vps
Environment=SMITHERS_RUNNER_TAGS=linux,node,git,shell,web,smithers
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env node /opt/runyard/src/runner.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
